# scripts/devbox-idle-stop.service
# /etc/systemd/system/devbox-idle-stop.service
[Unit]
Description=Devbox idle-stop check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/devbox-idle-stop