# scripts/devbox-idle-stop.timer
# /etc/systemd/system/devbox-idle-stop.timer
[Unit]
Description=Run devbox idle-stop check every 30 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target