# scripts/outpost-idle-stop.service
# /etc/systemd/system/outpost-idle-stop.service
[Unit]
Description=outpost idle-stop check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/outpost-idle-stop